# setups/rd-rclone/rclone-debrid-mount.service
[Unit]
Description=RClone Mount Debrid Service
After=network-online.target

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/rclone mount rd: /home/pi/rd --dir-cache-time 10s --multi-thread-streams=0 --cutoff-mode=cautious --vfs-cache-mode minimal --network-mode --buffer-size=0 --read-only --allow-other
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

